# src/rpm/scripts/postinstall
hasBeenEnabledOnStart=false

if [ $1 -eq 1 ] ; then 
    # Initial installation 
    if [ -x /bin/systemctl ] ; then
        /bin/systemctl daemon-reload >/dev/null 2>&1 || :
        /bin/systemctl enable elasticsearch.service
        /bin/systemctl start elasticsearch.service
        hasBeenEnabledOnStart=true
    fi

    if [ -x /sbin/chkconfig -a "$hasBeenEnabledOnStart" == "false" ] ; then
        /sbin/chkconfig --add elasticsearch
        # older suse linux distributions do not ship with systemd
        # but do not have an /etc/init.d/ directory
        # this tries to start elasticsearch on these as well without failing this script
        if [ -x /etc/init.d/elasticsearch ] ; then
            /etc/init.d/elasticsearch start
        elif [ -x /etc/rc.d/init.d/elasticsearch ] ; then
            /etc/rc.d/init.d/elasticsearch start
        fi
    fi

fi
